# deploy/upstart/clockwork.conf
# /etc/init/clockwork.conf - Clockwork config

description "Clockwork Background Worker"

start on runlevel [2345]
stop on runlevel [06]

# change to match your deployment user
setuid pi
setgid pi
env HOME=/home/pi

respawn
respawn limit 3 30

env RAILS_ENV=production

script
  cd /www/greateye/current/ && /home/pi/.rbenv/shims/bundle exec clockwork clock.rb
end script

